# roles/backup/files/backblaze_backup.timer
[Unit]
Description=Backblaze rclone timer

[Timer]
OnCalendar=3:0:0
Persistent=true

[Install]
WantedBy=timers.target
