# roles/backup/files/db_backup.timer
[Unit]
Description=Immich Postgres dump timer

[Timer]
OnCalendar=2:30:0
Persistent=true

[Install]
WantedBy=timers.target
