# roles/navidrome/files/nd_autoremove.timer
[Unit]
Description=Run navidrome autoremove

[Timer]
OnCalendar=*:30
Persistent=true

[Install]
WantedBy=timers.target
