# roles/gatus/templates/gatus-container-check.timer.j2
# {{ ansible_managed }} — roles/gatus.
[Unit]
Description=Run the Gatus container-state pusher periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec={{ gatus_push_interval_sec }}s
AccuracySec=15s

[Install]
WantedBy=timers.target
