# automation/roles/beelink/files/bw_backup.timer
[Unit]
Description=Run Bitwarden Backup Daily

[Timer]
OnCalendar=daily
Unit=bw_backup.service

[Install]
WantedBy=timers.target
