# automation/roles/navidrome/files/nd_remove_missing.timer
[Unit]
Description=Run ND remove missing

[Timer]
OnCalendar=*:25
Persistent=true

[Install]
WantedBy=timers.target
